# config-linux/.config/containers/systemd/open-webui.container
[Container]
Image=ghcr.io/open-webui/open-webui:main
AddDevice=nvidia.com/gpu=all
AutoUpdate=registry
Volume=open-webui:/app/backend/data
PublishPort=6002:8080
SecurityLabelDisable=true

[Service]
ExecStartPre=-sudo nvidia-ctk cdi generate --output=/etc/cdi/nvidia.yaml
Restart=always

[Install]
WantedBy=default.target
